# cron.d/patch_tuesday
0 23 * * 2 root [ $(date +\%d) -ge 7 ] && [ $(date +\%d) -le 14 ] && /opt/scripts/patch_tuesday.sh
